# dots/systemd/user/calcurse-caldav.service
[Unit]
Description=Sync calcurse calendar

[Service]
Type=oneshot
ExecStart=calcurse-caldav
